# jasper-watchdog-v2/tests/fixtures/crontab
#!/usr/bin/env bash
echo "0 13 * * * /home/opc/script/limpiar_tomcat.sh"
